# contrib/init/ruxcoind.conf
description "Ruxcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env RUXCOIND_BIN="/usr/bin/ruxcoind"
env RUXCOIND_USER="ruxcoin"
env RUXCOIND_GROUP="ruxcoin"
env RUXCOIND_PIDDIR="/var/run/ruxcoind"
# upstart can't handle variables constructed with other variables
env RUXCOIND_PIDFILE="/var/run/ruxcoind/ruxcoind.pid"
env RUXCOIND_CONFIGFILE="/etc/ruxcoin/ruxcoin.conf"
env RUXCOIND_DATADIR="/var/lib/ruxcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # ruxcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$RUXCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run ruxcoind."
        echo "The setting must appear in $RUXCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Ruxcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$RUXCOIND_PIDDIR"
    chmod 0755 "$RUXCOIND_PIDDIR"
    chown $RUXCOIND_USER:$RUXCOIND_GROUP "$RUXCOIND_PIDDIR"
    chown $RUXCOIND_USER:$RUXCOIND_GROUP "$RUXCOIND_CONFIGFILE"
    chmod 0660 "$RUXCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$RUXCOIND_PIDFILE" \
    --chuid $RUXCOIND_USER:$RUXCOIND_GROUP \
    --exec "$RUXCOIND_BIN" \
    -- \
    -pid="$RUXCOIND_PIDFILE" \
    -conf="$RUXCOIND_CONFIGFILE" \
    -datadir="$RUXCOIND_DATADIR" \
    -disablewallet \
    -daemon

